# .local/share/external/services/pulseaudio.service
[Unit]
Description=PulseAudio system server

[Service]
Type=notify
ExecStart=pulseaudio --daemonize=no --system --realtime --log-target=journal --disallow-exit --disallow-module-loading

[Install]
WantedBy=multi-user.target
